# crontab
30 3 * * * /path/to/venv/bin/python crawl.py --page_id=cocacola